# obsidian/configs/obsidian-mcp-go.service
[Unit]
Description=Obsidian Custom MCP Server (Go)
After=obsidian-headless.service

[Service]
Type=simple
EnvironmentFile=/etc/obsidian-mcp.env
ExecStart=/usr/local/bin/obsidian-mcp /opt/obsidian-vault 3333
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
